# xyz.endymion.avella-tray.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Avella tray (menu bar) launchd agent.

  Install:
    cp xyz.endymion.avella-tray.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/xyz.endymion.avella-tray.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/xyz.endymion.avella-tray.plist
    rm ~/Library/LaunchAgents/xyz.endymion.avella-tray.plist

  Launches via `open -W` so LaunchServices associates the process with the
  .app bundle (required for notification icons and proper app identity).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>xyz.endymion.avella-tray</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/open</string>
        <string>-W</string>
        <string>/Applications/AvellaTray.app</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/shrike/Library/Logs/avella-tray.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/shrike/Library/Logs/avella-tray.log</string>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
